# 99-mountcopy.rules
ACTION=="add", SUBSYSTEMS=="usb", KERNEL=="sd*", RUN+="/bin/sh -c '/bin/sleep 2 ; /usr/bin/mountcopy %k'"
